# etc/profile.d/npm-aliases.sh
alias n='npm'
alias nb='npm build'
alias nc='npm cache'
alias nca='npm cache add'
alias ncc='npm cache clean'
alias ncv='npm cache verify'
alias nh='npm help'
alias ni="npm install"
alias nid='npm install --save-dev'
alias nig="npm install -g"
alias nin="npm init"
alias nli="npm login"
alias nlo="npm logout"
alias np="npm publish"
alias npc="npm pack"
alias nr='npm run'
alias nrm="npm uninstall"
alias nu="npm update"
alias nug="npm update -g"
alias nup="npm unpublish"

# Vue, Nuxt, React, Next
alias nrb='npm run build'

# Vue
alias nrl='npm run lint'
alias nrsr='npm run serve'
alias nrt='npm run test'

# Nuxt
alias nrd='npm run dev'
alias nrg='npm run generate'
alias nrs='npm run start'

# Next
alias ns='npm start'
alias nt='npm test'

# Laravel Mix
alias nrh='npm run hot'
alias nrp='npm run production'
alias nrw='npm run watch'
alias nrwp='npm run watch-poll'
