# etc/profile.d/npm-aliases.sh
alias n='npm'
alias nh='npm help'
alias nid='npm i --save-dev'
alias nig="npm i -g"
alias nin="npm init"
alias nli="npm login"
alias nlo="npm logout"
alias np="npm publish"
alias npc="npm pack"
alias nr='npm run'
alias nug="npm up -g"

# Vue, Nuxt, React, Next
alias nrb='npm run build'

# Vue
alias nrl='npm run lint'
alias nrsr='npm run serve'
alias nrt='npm run test'

# Nuxt
alias nrd='npm run dev'
alias nrg='npm run generate'
alias nrst='npm run start'

# Next
alias ns='npm start'
alias nt='npm test'
